# net-fs/samba-tng/files/init/samba-tng
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /home/cvsroot/gentoo-x86/net-fs/samba/files/samba-init,v 1.1 2002/05/06 19:57:08 woodchip Exp $
#
# Start/stop script for Samba TNG
#
# This is from a Solaris box.  Proper links need to be made to it from the
# rcX.d directories.  For Linux, it goes into /etc/rc.d/init.d
# 
# Lonnie
#
# Heavily modified and integrated with autoconf by Elrond
#
# RedHat 7.0 or above puts this in /etc/init.d
# chkconfig: - 91 35
# description: Starts and stops the Samba TNG daemons

prefix=/usr
exec_prefix=${prefix}

bindir=/usr/sbin
sbindir=/usr/sbin
localstatedir=/var/log
tmpdir=/var/run/samba-tng
piddir=$tmpdir
sysconfdir=/etc/samba-tng
etcdir=/etc/samba-tng

daemons="nmbd smbd"
daemons="$daemons
  netlogond
  samrd
  browserd
  lsarpcd
  srvsvcd
  winregd
  wkssvcd
  spoolssd
  svcctld
"
unset TMPDIR

depend() 
{
	need net
	use cupsd
}

start() 
{
	ebegin "Starting samba-TNG"
	for i in $daemons
	do
		echo -ne "    Starting $i\n"
		start-stop-daemon --start --quiet --exec /usr/sbin/${i}
		result=$?
	done
	result=$(( $result + $? ))
	eend $result
}

stop() 
{
	ebegin "Stopping samba-TNG"
	for i in $daemons
	do
		echo -ne "    Stopping $i\n"
		start-stop-daemon --stop --quiet --pidfile /var/run/samba-tng/${i}.pid
		result=$?
	done
	result=$(( $result + $? ))
	eend $result
}




