# net-fs/samba.my/files/init/daemontools/majorversion/3.0/init.d/svc-samba
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: $

SVC_DIR=/var/lib/supervise
service_name=samba

depend () 
{
    need net
    after net
}

start () 
{
    ebegin "Starting ${service_name}"
    if [ ! -e ${SVC_DIR}/control/smbd ] ; 
	then
	ln -sf /service/smbd ${SVC_DIR}/control/smbd
    fi
    if [ ! -e ${SVC_DIR}/control/nmbd ] ; 
	then
	ln -sf /service/nmbd ${SVC_DIR}/control/nmbd
    fi
    /usr/bin/svc -u ${SVC_DIR}/control/smbd
    /usr/bin/svc -u ${SVC_DIR}/control/nmbd
    eend $?
}

stop () 
{
    ebegin "Stopping ${service_name}"
    if [ -e ${SVC_DIR}/control/smbd ] ; 
	then
	/usr/bin/svc -dx ${SVC_DIR}/control/smbd/log
	/usr/bin/svc -dx ${SVC_DIR}/control/smbd
	rm ${SVC_DIR}/control/smbd
    fi
    if [ -e ${SVC_DIR}/control/nmbd ] ; 
	then
	/usr/bin/svc -dx ${SVC_DIR}/control/nmbd/log
	/usr/bin/svc -dx ${SVC_DIR}/control/nmbd
	rm ${SVC_DIR}/control/nmbd
    fi
    eend $?
}


# Local Variables:
# mode: sh
# End:
