# gnustep-base/gnustep-base-/files/extra/init/gentoo-sysv/majorversion/1/init.d/gnustep
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: $

if [ -f /etc/env.d/10gnustep_root ]
    then
    sh /etc/env.d/10gnustep_root
fi


if [ -z "${GNUSTEP_ROOT}" ]
    then
    GNUSTEP_ROOT=/usr/GNUstep
fi

source ${GNUSTEP_ROOT}/System/Library/Makefiles/GNUstep.sh

if [ -z "$GNUSTEP_FLATTENED" ]; 
    then
    TDIR=${GNUSTEP_SYSTEM_ROOT}/Tools/${GNUSTEP_HOST_CPU}/${GNUSTEP_HOST_OS}
else
    TDIR=${GNUSTEP_SYSTEM_ROOT}/Tools
fi


GDOMAP=${TDIR}/gdomap
GDNC=${TDIR}/${LIBRARY_COMBO}/gdnc

checkfiles () 
{
    # Check for necessary files
    
    if [ ! -d ${GNUSTEP_ROOT}/System/Makefiles -o ! -x ${GDOMAP} -o ! -x ${GDNC} ] 
	then
	eerror "GNUstep is not installed"
	return 1
    fi
    return 0
}


start () 
{
    local retval
    local retval1
    local retval2
    
    checkfiles || return $?

    ebegin "Starting GNUstep services"
    
    GNUSTEP_MAKEFILES="${GNUSTEP_ROOT}/System/Makefiles"
    
    # Source GNUstep makefile
    if [ -x $GNUSTEP_MAKEFILES/GNUstep.sh ] 
	then
	. $GNUSTEP_MAKEFILES/GNUstep.sh
    else
	exit 1
    fi
    
    # Start services

    start-stop-daemon --start --pidfile /var/run/gdomap.pid --exec ${GDOMAP}  -- -I /var/run/gdomap.pid
    retval1=$?
    start-stop-daemon --start --exec ${GDNC}
    retval2=$?
    retval=$((retval1 + retval2))
    eend ${retval}
}

stop () 
{
    local retval
    local retval1
    local retval2
    
    ebegin "Stopping GNUstep services"
    
    start-stop-daemon --stop --exec ${GDNC}
    retval1=$?
    start-stop-daemon --stop --pidfile /var/run/gdomap.pid --exec ${GDOMAP}
    rm /var/run/gdomap.pid
    retval2=$?
    retval=$((retval1 + retval2))
    eend ${retval}
}



# Local Variables:
# mode: sh
# End:

