# net-firewall/pcx_firewall/files/rc.d/2.21/pcx_firewall
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
#
# Startup script to run the PCX Firewalls generated scripts.
#
# chkconfig: 2345 11 92
#
# description: Automates a packet filtering firewall with iptables.

LOCATION="/etc/pcx_firewall"

opts="${opts} clear panic"

depend() {
	need net
}

checkconfig() {
	# make sure we have the necessary scripts
	if [ ! -e $LOCATION/firewall ]; then
	    eerror "Check for firewall script: \$LOCATION/firewall"
	    return 1	    
	fi
	return 0
}

start() {
	ebegin "Starting pcx_firewall"
	$LOCATION/firewall start
	eend $?
}

stop() 
{
	ebegin "Stopping pcx_firewall"
	$LOCATION/firewall stop
	eend $?
}

restart() 
{
	ebegin "Restarting pcx_firewall"
	$LOCATION/firewall restart
	eend $?
}

status() 
{
  $LOCATION/firewall status 
}

panic()
{
  $LOCATION/firewall panic 
}

clear() 
{ 
  ebegin "Clear pcx_firewall"
  $LOCATION/firewall clear 
  eend $?
}
