# mail-filter/sa-exim/files/sa-exim_greylistclean.cron
#!/bin/sh

[ -x /usr/sbin/greylistclean ] && /usr/sbin/greylistclean
